# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

# 示例：还可以放一个 web 服务
#[program:myapp]
#command=python3 /app/server.py
#autorestart=true

